# config/flask_cbr.conf
[program:flask_cbr]
command = gunicorn flask_cbr:app -b localhost:8000
directory = /home/flask_cbr/flask_cbr/
user = root